# etc/inventory-repair-app.service
# cp /etc/systemd/system/inventory-repair-app.service
# systemctl daemon-reload
# systemctl start inventory-repair-app.service
# systemctl status inventory-repair-app.service
#
# To enable by default on boot:
#
# systemctl enable inventory-repair-app.service

[Unit]
Description=Yesdo inventory server
After=network.target
StartLimitIntervalSec=0

[Service]
WorkingDirectory=/home/amauricio/inventory-repair-app
Type=simple
User=amauricio
Restart=always
RestartSec=1
ExecStart=/home/amauricio/inventory-repair-app/.stack-work/install/i386-linux/082105eb77d1e1b3b5ed46c9697f28adfb650dff873ea261608b680fdd548c1a/8.6.5/bin/inventory-repair-app

[Install]
WantedBy=multi-user.target
